# vm/storage/rootfs/etc/systemd/system/postgresql@14-protect.service.d/10-assert-datadir.conf
# Assert data_directory=/data before postgres starts. See
# /usr/local/sbin/assert-protect-datadir.sh for the full rationale.
[Service]
ExecStartPre=/usr/local/sbin/assert-protect-datadir.sh
