# host/com.protect-on-mac.vm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.protect-on-mac.vm.plist

  Launches the Protect VM at boot and keeps it running. Place this file
  at /Library/LaunchDaemons/ and load it with launchctl. See the
  install-launchd.sh script in this repo for a one-shot installer.

  TEMPLATE NOTES:
    Replace YOUR_USERNAME with the macOS user the script runs as. This
    user needs the passwordless qemu sudoers entry from the README
    setup section. Replace SCRIPT_PATH with the absolute path to your
    start-protect-vm.sh.

  WHY LaunchDaemon (not LaunchAgent):
    LaunchDaemons run as root at boot, even before any user logs in.
    LaunchAgents only run inside a logged-in user session. For a
    headless server you want LaunchDaemon — the VM should start whether
    or not anyone has logged in to the Mac console.

  WHY KeepAlive = SuccessfulExit false:
    start-protect-vm.sh translates the VM's fate into its own exit code.
    A guest reboot is handled inside the script — it cold-restarts QEMU
    itself and never exits. A deliberate power-off (e.g. "Power Off" from
    the Protect web UI) exits 0; SuccessfulExit=false tells launchd NOT
    to restart on a clean exit, so the VM stays down until someone starts
    it again. A crash or QEMU fault exits non-zero, and launchd does
    relaunch. The 30-second throttle prevents tight restart loops if
    something is consistently failing.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.protect-on-mac.vm</string>

    <key>UserName</key>
    <string>YOUR_USERNAME</string>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME</string>

    <!-- The actual command to run. ProgramArguments wins over Program
         when both are present, and ProgramArguments lets us pass args
         explicitly even if there aren't any here. -->
    <key>ProgramArguments</key>
    <array>
        <string>SCRIPT_PATH</string>
    </array>

    <!-- Start at boot. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart only on an abnormal (non-zero) exit. A clean exit 0 means
         the VM was deliberately powered off — leave it down. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Don't restart faster than every 30 seconds. Prevents busy-loops
         if QEMU is consistently failing (e.g. DAS not connected). -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Capture pre-QEMU script output to a log file. This is the output
         from start-protect-vm.sh itself (disk resolution, mode detection,
         "Starting VM..." messages, etc.) plus anything QEMU writes to
         stderr.

         The VM's actual console output (kernel boot, login prompt, etc.)
         is captured separately via the chardev logfile inside QEMU, in
         /var/log/protect-vm.console.log. This is so you can attach to the
         live console via attach-console.sh AND have a persistent log of
         everything that's happened on the console.
    -->
    <key>StandardOutPath</key>
    <string>/var/log/protect-vm.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/protect-vm.error.log</string>

    <!-- A sane PATH (launchd does not inherit a shell's), plus the path to
         THIS VM's protect-on-mac.conf. The conf lives in the VM's data
         directory, so start-protect-vm.sh — started here by launchd, with
         no shell environment — is told explicitly where to find it. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PROTECT_ON_MAC_CONF</key>
        <string>CONF_PATH</string>
    </dict>
</dict>
</plist>
