# vm/storage/rootfs/etc/systemd/system/md-health-watch.service
[Unit]
Description=md array health watcher (Pushover on degraded/failed)
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/md-health-watch.sh
